# systemd/qbt.container
[Unit]
Description=qBittorrent (static, rootless, foreground)
Wants=network-online.target
After=network-online.target

[Container]
# Image built locally as qbt:latest
Image=localhost/qbt:latest
Network=pasta
DNS=192.168.2.2

# Rootless container
UserNS=keep-id

# Bind mounts (host -> container)
Volume=/etc/localtime:/etc/localtime:ro
Volume=/srv/containers/qbt/config:/qbt/config:rw
Volume=/srv/containers/qbt/incomplete:/qbt/incomplete:rw
Volume=/srv/containers/qbt/downloads:/qbt/downloads:rw
Volume=/srv/containers/qbt/monitor:/qbt/monitor:rw

# Publish ports
# WebUI (only one exposed in the image)
PublishPort=1340:1340/tcp

# BitTorrent ports (TCP + UDP)
PublishPort=8999:8999/tcp
PublishPort=8999:8999/udp

# Embedded tracker
PublishPort=9000:9000/tcp

# Run in foreground (ENTRYPOINT already set)
# No need for Exec= override
[Service]
Type=simple
NotifyAccess=none
# Restart policy
Restart=always
RestartSec=36
TimeoutStopSec=15

[Install]
WantedBy=default.target
