# systemd/user/lektrix-trend-cp.timer
# This timer is for triggering the trending of data

[Unit]
Description=periodic scp & trending electricity data (timer)

[Timer]
OnCalendar=hourly
RandomizedDelaySec=1m

[Install]
WantedBy=timers.target
